# amdgpufan.service
[Unit]
Description=AMDGPU fan
After=multi-user.target

[Service]
ExecStart=/usr/bin/amdgpufan

[Install]
WantedBy=multi-user.target
